# crontab
*/10 * * * * root /pancache-sync /pancache/peers.txt
